# demos/etcd/etcd-node.service
[Unit]
Description=etcd cluster node (qarax demo)
Documentation=https://etcd.io
# The kernel ip= param configures eth0 before systemd starts, so we only need
# the basic network target — no need for network-online.target.
After=network.target

[Service]
Type=exec
ExecStart=/usr/local/bin/etcd-node.sh
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
